# scripts/claude-puppet-dashboard.service
[Unit]
Description=claude-puppet dashboard (mission-control UI on :5055)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/projects/claude-puppet/dashboard
# Use the user's nvm-managed node so PATH resolves npm/node correctly.
Environment=PATH=%h/.nvm/versions/node/v24.14.1/bin:/usr/local/bin:/usr/bin:/bin
Environment=NODE_ENV=production
Environment=PORT=5055
ExecStart=%h/.nvm/versions/node/v24.14.1/bin/npm start
Restart=on-failure
RestartSec=2s
StandardOutput=journal
StandardError=journal
# Cap memory + cpu so a runaway dashboard can't starve the box.
MemoryMax=512M
CPUQuota=80%

[Install]
WantedBy=default.target
